# Library/LaunchAgents/offlineimap.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>offlineimap</string>
    <key>ProgramArguments</key>
    <array>
      <string>/Users/bert/bin/offlineimap</string>
    </array>
    <key>RunAtLoad</key>
    <false/>
    <key>EnableTransactions</key>
    <false/>
    <key>ServiceDescription</key>
    <string>OfflineIMAP</string>
    <key>UserName</key>
    <string>bert</string>
    <key>StandardErrorPath</key>
    <string>/Users/bert/.offlineimap/err.log</string>
    <key>StandardOutPath</key>
    <string>/Users/bert/.offlineimap/out.log</string>
    <key>EnvironmentVariables</key>
    <dict>
      <key>PATH</key>
      <string>/Users/bert/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
      <key>OFFLINEIMAP_RUN_FROM_LAUNCHD</key>
      <string>true</string>
    </dict>
  </dict>
</plist>
